# deploy/latent_self.service
[Unit]
# Installed by deploy/install_kiosk.sh.
# The User= value below must match the KIOSK_USER variable in that script.
Description=Latent Self Kiosk
After=network.target

[Service]
ExecStart=/opt/latent-self/latent-self --ui qt
Restart=on-failure
WatchdogSec=30s
Type=simple
User=kiosk

[Install]
WantedBy=multi-user.target
